# deploy/systemd/amm-web.service
[Unit]
Description=AMM Web (Reflex)
After=network.target amm-core.service
Wants=amm-core.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=pegasus
Group=pegasus
WorkingDirectory=/home/pegasus/AMM_web_prod

# Ensure reflex can spawn granian in prod single-port mode.
Environment=PATH=/home/pegasus/AMM_web_prod/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=REFLEX_DIR=/home/pegasus/AMM_web_prod/.reflex
Environment=REFLEX_USE_NPM=1
Environment=AMM_REQUIRED_NODE_VERSION=20.19.0

ExecStart=/home/pegasus/AMM_web_prod/scripts/run_single_port.sh

Restart=on-failure
RestartPreventExitStatus=78
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
